# src/assets/admin-notify/units/admin-checkupdates.service
[Unit]
Description=Collect pending package updates summary for admin notifications

[Service]
Type=oneshot
ExecStart=/usr/local/sbin/admin-checkupdates-collect
